# .config/systemd/user/pkgs.timer
[Unit]
Description=Shows number of package updates available

[Timer]
OnBootSec=2m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
